# profiles/tablet-niri/rootfs-overlay/etc/systemd/system/tb321fu-usb-rescue.service
[Unit]
Description=TB321FU USB network and serial rescue gadget
Wants=NetworkManager.service sys-kernel-config.mount
After=NetworkManager.service systemd-modules-load.service sys-kernel-config.mount

[Service]
Type=simple
ExecStart=/usr/local/libexec/tb321fu-usb-rescue
Restart=always
RestartSec=5s
TimeoutStopSec=15s

[Install]
WantedBy=multi-user.target
