# profiles/tablet-niri/rootfs-overlay/etc/systemd/user/noctalia.service
[Unit]
Description=Noctalia desktop shell
Documentation=https://github.com/noctalia-dev/noctalia
PartOf=graphical-session.target
After=graphical-session.target pipewire.service
Requisite=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/noctalia
Restart=on-failure
RestartSec=2s
TimeoutStopSec=8s

[Install]
WantedBy=niri.service
