# source/tb321fu-camera-rootfs-overlay/rootfs-overlay/etc/systemd/user/pipewire.service.d/50-y700-libcamera-ipa.conf
[Service]
# libcamera's isolated IPA proxy uses clone3() with CLONE_NEWUSER and
# CLONE_NEWNET. Ubuntu's default PipeWire sandbox blocks namespace creation,
# which prevents the Y700 source-built SoftISP from running through PipeWire.
# SoftIPA needs only user and network namespaces; retain all other namespace
# restrictions from the PipeWire service sandbox.
RestrictNamespaces=user net
